# config/supervisor/ddns.conf
[program:ddns]
command=/usr/local/bin/node /var/lib/nsd/ddns/server.js
directory=/var/lib/nsd
user=nsd
autostart=true
autorestart=true
redirect_stderr=true
