# config/supervisor/etherpad.conf
[program:etherpad]
command=/var/www/etherpad/bin/run.sh
directory=/var/www/etherpad
user=etherpad
autostart=true
autorestart=true
